# etc/supervisor/conf.d/datanode.conf.erb
[program:datanode]
user=hadoop
command=/opt/hadoop-2.7.0/bin/hdfs datanode