# etc/supervisor/conf.d/datastax-agent.conf.erb
[program:datastax-agent]
user=cassandra
command=/opt/datastax-agent-5.2.1/bin/datastax-agent -f
